# baseconf/mysql/SYSTEM/mysqlloglotate.service
[Unit]
Description=mysql log lotate
Requires=docker.service

[Service]
ExecStart=/usr/bin/docker exec conkan /etc/cron.daily/logrotate
Type=oneshot
