# baseconf/conkan/SYSTEM/conkanlogrotate.service
[Unit]
Description=conkan log rotate
Requires=docker.service

[Service]
ExecStart=/usr/bin/docker exec conkan /etc/cron.daily/logrotate
Type=oneshot
